# deployment/pixelsort-kiosk.service
[Unit]
Description=Raspberry Pi Pixel Sorter (Auto-updating Kiosk Mode)
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pixelsort
WorkingDirectory=/home/pixelsort/Pixelsort
ExecStart=/home/pixelsort/Pixelsort/deployment/run_pixelsort.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables for X11 display
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pixelsort/.Xauthority"

[Install]
WantedBy=graphical.target
